# services/unified-guide/unified-guide-refresh.timer
[Unit]
Description=Periodic Unified Guide refresh (every 15 min)

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
