# services/yttv-scraper/yttv-scraper.timer
[Unit]
Description=Periodic YouTube TV scrape (every 30 min)

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
